# files/launch/ls-hubd_private.conf.in
# @@@LICENSE
#
#      Copyright (c) 2012-2014 LG Electronics, Inc.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
# http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
# LICENSE@@@

description "Luna Private Hub service"

start on stopped bootmisc
stop on started start_update
emits ls-hubd_private-ready

# OOM should never kill ls-hubs
oom score never

# Comment this line out to suppress logs on the console
#console output

pre-start script
	mkdir -p @WEBOS_INSTALL_WEBOS_LOCALSTATEDIR@/ls2/roles/prv @WEBOS_INSTALL_WEBOS_LOCALSTATEDIR@/ls2/roles/pub @WEBOS_INSTALL_WEBOS_LOCALSTATEDIR@/ls2/services/prv @WEBOS_INSTALL_WEBOS_LOCALSTATEDIR@/ls2/services/pub
end script

script
	# Use centralized nice script
	if [ -f /etc/init.d/webos-nicer.sh ]
	then
		. /etc/init.d/webos-nicer.sh
	fi
	exec $WEBOS_NICE @WEBOS_INSTALL_SBINDIR@/ls-hubd --conf @WEBOS_INSTALL_SYSCONFDIR@/luna-service2/ls-private.conf
end script
